# supervisord/conf.d/mockoon.conf
[program:mockoon]
command=mockoon --no-sandbox --disable-gpu
